# systemd/bluetooth-reconnectd.service
[Unit]
Description=Bluetooth auto-reconnect daemon
After=bluetooth.target

[Service]
ExecStart=/usr/local/bin/bluetooth-reconnectd
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
